# files/etc/init.d/custom-model
#!/bin/sh /etc/rc.common

START=95
STOP=10

normalize() {
    local value="$1"

    value="$(printf '%s' "$value" | tr '\t\r\n' '   ' | sed 's/[[:space:]][[:space:]]*/ /g; s/^ //; s/ $//')"
    printf '%s' "$value"
}

read_dmi() {
    local path="$1"

    [ -r "$path" ] || return 0
    tr -d '\000' < "$path" 2>/dev/null
}

is_placeholder() {
    case "$1" in
        ""|"Default string"|"To be filled by O.E.M."|"To Be Filled By O.E.M."|"Not Specified"|"Not Applicable"|"System Product Name"|"System Version"|"Undefined"|\
        "OpenWrt"|"Generic"|"Generic PC"|"Standard PC (i440FX + PIIX, 1996)"|"Virtual Machine"|"KVM"|"Bochs"|"QEMU"|"OpenStack Foundation"|"Microsoft Corporation Virtual Machine")
            return 0
            ;;
    esac

    return 1
}

detect_virtual_model() {
    local vendor="$1"
    local product="$2"
    local version="$3"
    local combined

    combined="$(normalize "$vendor $product $version")"

    case "$combined" in
        *QEMU*|*qemu*|*Standard\ PC*|*Bochs*)
            printf '%s\n' "QEMU Virtual Router"
            return 0
            ;;
        *KVM*|*kvm*)
            printf '%s\n' "KVM Virtual Router"
            return 0
            ;;
        *VMware*|*vmware*)
            printf '%s\n' "VMware Virtual Router"
            return 0
            ;;
        *VirtualBox*|*virtualbox*|*innotek*)
            printf '%s\n' "VirtualBox Router"
            return 0
            ;;
        *Hyper-V*|*Microsoft\ Corporation\ Virtual\ Machine*)
            printf '%s\n' "Hyper-V Virtual Router"
            return 0
            ;;
        *OpenStack*)
            printf '%s\n' "OpenStack Virtual Router"
            return 0
            ;;
    esac

    return 1
}

pick_model() {
    local current vendor product version candidate

    current="$(normalize "$(cat /tmp/sysinfo/model 2>/dev/null)")"
    if [ -n "$current" ] && ! is_placeholder "$current"; then
        return 0
    fi

    vendor="$(normalize "$(read_dmi /sys/class/dmi/id/sys_vendor)")"
    product="$(normalize "$(read_dmi /sys/class/dmi/id/product_name)")"
    version="$(normalize "$(read_dmi /sys/class/dmi/id/product_version)")"

    candidate="$(detect_virtual_model "$vendor" "$product" "$version")" || true
    if [ -n "$candidate" ]; then
        printf '%s\n' "$candidate"
        return 0
    fi

    if ! is_placeholder "$product"; then
        if ! is_placeholder "$vendor"; then
            printf '%s\n' "$(normalize "$vendor $product")"
        else
            printf '%s\n' "$product"
        fi
        return 0
    fi

    if ! is_placeholder "$vendor"; then
        printf '%s\n' "$vendor Router"
        return 0
    fi

    printf '%s\n' "x86-64 OpenWrt Router"
}

apply_model() {
    local model

    mkdir -p /tmp/sysinfo
    model="$(pick_model)"
    [ -n "$model" ] || return 0
    printf '%s\n' "$model" > /tmp/sysinfo/model
}

boot() {
    apply_model
}

start() {
    apply_model
}
