# files/etc/init.d/owb-boot
#!/bin/sh /etc/rc.common

START=99
STOP=01

run_builder_boot() {
    if grep -qiF "Default string" /sys/class/dmi/id/product_name 2>/dev/null; then
        mkdir -p /tmp/sysinfo
        echo "OpenWrt Router" > /tmp/sysinfo/model
    fi

    if [ -x /usr/bin/zb ]; then
        /usr/bin/zb >/dev/null 2>&1 || true
    fi
}

boot() {
    start "$@"
}

start() {
    run_builder_boot
}
