# files/etc/rc.local
#!/bin/sh

if grep -qiF "Default string" /sys/class/dmi/id/product_name 2>/dev/null; then
    mkdir -p /tmp/sysinfo
    echo "OpenWrt Router" > /tmp/sysinfo/model
fi

exit 0
